# files/supervisord.conf
[supervisord]
nodaemon=true

[program:aria2]
command=aria2c --conf-path=/conf/aria2.conf

[program:webui]
command=darkhttpd /aria2-webui